# rsnapshot@.timer
[Unit]
Description=rsnapshot %i backup

[Timer]
OnCalendar=%i
Persistent=true
Unit=rsnapshot-%i.service

[Install]
WantedBy=timers.target
